# .aws/crontasks/timers/send_monthly_manager_pr_notices.timer
[Unit]
Description="Run send_monthly_manager_pr_notices.service on the first of the month at 8 AM"

[Timer]
OnCalendar=*-*-01 08:00:00 America/Los_Angeles
# Add a randomized activation delay of up to 1 minute to avoid activating
# multiple timers at exactly the same time
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target