# .aws/crontasks/timers/send_is_timeoff_next_week_report.timer
[Unit]
Description="Run send_is_timeoff_next_week_report.service every Friday at 8 AM"

[Timer]
OnCalendar=Fri *-*-* 08:00:00 America/Los_Angeles
# Add a randomized activation delay of up to 1 minute to avoid activating
# multiple timers at exactly the same time
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target